# scripts/podman/flowhelm.container.in
# FlowHelm gateway — Podman Quadlet (rootless)
# Installed by setup-podman.sh into flowhelm's ~/.config/containers/systemd/
# {{FLOWHELM_HOME}} is replaced at install time.

[Unit]
Description=FlowHelm gateway (rootless Podman)

[Container]
Image=flowhelm:local
ContainerName=flowhelm
UserNS=keep-id
Volume={{FLOWHELM_HOME}}/.flowhelm:/home/node/.flowhelm
EnvironmentFile={{FLOWHELM_HOME}}/.flowhelm/.env
Environment=HOME=/home/node
Environment=TERM=xterm-256color
PublishPort=18789:18789
PublishPort=18790:18790
Pull=never
Exec=node dist/index.js gateway --bind lan --port 18789

[Service]
TimeoutStartSec=300
Restart=on-failure

[Install]
WantedBy=default.target
